# roles/jboss/templates/jboss.service.j2
[Unit]
Description=JBoss WildFly Application Server
After=network.target

[Service]
Type=simple
User=jboss
Group=jboss
# On supprime le security manager qui bloque tout
# On utilise --server-config=standalone.xml pour être explicite
ExecStart=/opt/jboss/bin/standalone.sh -b 0.0.0.0 -bmanagement 0.0.0.0 --server-config=standalone.xml
# On utilise le script de shutdown standard
ExecStop=/opt/jboss/bin/jboss-cli.sh --connect command=:shutdown
# On simplifie le démarrage
Restart=on-failure
RestartSec=10
WorkingDirectory=/opt/jboss

[Install]
WantedBy=multi-user.target